# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid frank_berghaus
chdir /home/frank_berghaus/sites/SITENAME/source
exec /home/frank_berghaus/.virtualenvs/virtualenv/bin/gunicorn \
     --bind unix:/tmp/SITENAME.socket \
     superlists.wsgi:application
